# lesson6_SystemD/SCRIPTS/apache_multipleConf/configs/httpd.target
[Unit]
Wants=httpd@site1.service httpd@site2.service
